# s6-overlay/s6-rc.d/nsd/finish
#!/usr/bin/execlineb -S0

if { eltest $# -ne 0 }
if { eltest ${1} -ne 256 }
s6-svscanctl -t /var/run/service