# s6-overlay/s6-rc.d/nsd-init/up
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

foreground { echo "Initializing temporary nsd working directories" }
/etc/s6-overlay/scripts/init.sh